# openepm-agent.service
[Unit]
Description=OpenEPM Agent
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=/home/handesh/Documents/Uni/project/agent
Environment="PATH=/home/handesh/Documents/Uni/project/agent/.venv/bin"
ExecStart=/home/handesh/Documents/Uni/project/agent/.venv/bin/python -m openepm_agent.cli start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
